# resources/crontab
# sudo
@reboot python3.10 ~/ArchieHub/SmartHome/hardware/WiFi.py
